# docker-lesson2/node/container_files/etc/supervisord.d/node.ini
[program:node]
directory = /app
command = node app.js
